# launchagents/com.arkadiy.smart-home-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.arkadiy.smart-home-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/arkadiysudarikov/Library/Application Support/SmartHomeMonitor' &amp;&amp; ./scripts/smart_home_snapshot.py &amp;&amp; ./scripts/maintain_storage.py &amp;&amp; ./scripts/analyze_patterns.py &amp;&amp; ./scripts/analyze_energy_pairing.py &amp;&amp; ./scripts/analyze_all_energy_readings.py &amp;&amp; ./scripts/capture_alarm_com.js &amp;&amp; ./scripts/analyze_chargepoint_pairing.py &amp;&amp; ./scripts/analyze_meter_reconciliation.py &amp;&amp; ./scripts/analyze_bill_home_pairing.py &amp;&amp; ./scripts/analyze_energy_costs.py &amp;&amp; ./scripts/analyze_combined_energy_monitor.py &amp;&amp; ./scripts/generate_alerts.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/arkadiysudarikov/Library/Application Support/SmartHomeMonitor</string>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/arkadiysudarikov/Library/Application Support/SmartHomeMonitor/logs/monitor.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/arkadiysudarikov/Library/Application Support/SmartHomeMonitor/logs/monitor.err.log</string>
</dict>
</plist>
